# ansible/templates/treebeard-router.service.j2
[Unit]
Description=Treebeard Router
Requires=network.target
After=syslog.target network.target

[Service]
Type=simple
ExecStart=/root/treebeard/router/router -routerid {{ item.id }} -ip {{ item.local_bind_ip }} -port {{ item.port }} -conf /root/treebeard/ -logpath /root/treebeard/router/{{ item.id }}/all.log
Restart=on-failure

[Install]
WantedBy=multi-user.target